# packaging/macos/pkg/scripts/postinstall
#!/bin/bash
set -eu

# Allocate UID/GID for the _aiproxy service account if absent.
# Pick the lowest free UID in the 220-400 macOS system-user range.
if ! dscl . -read /Users/_aiproxy UniqueID >/dev/null 2>&1; then
  uid=""
  for candidate in $(seq 220 400); do
    if ! dscl . -search /Users UniqueID "$candidate" 2>/dev/null | grep -q .; then
      uid="$candidate"
      break
    fi
  done
  if [ -z "$uid" ]; then
    echo "ai-proxy postinstall: no free UID in 220-400; allocate _aiproxy manually." >&2
    exit 1
  fi

  dscl . -create /Groups/_aiproxy
  dscl . -create /Groups/_aiproxy PrimaryGroupID "$uid"
  dscl . -create /Groups/_aiproxy RealName "AI Proxy Service"
  dscl . -create /Groups/_aiproxy Password "*"

  dscl . -create /Users/_aiproxy
  dscl . -create /Users/_aiproxy UniqueID "$uid"
  dscl . -create /Users/_aiproxy PrimaryGroupID "$uid"
  dscl . -create /Users/_aiproxy NFSHomeDirectory /var/empty
  dscl . -create /Users/_aiproxy UserShell /usr/bin/false
  dscl . -create /Users/_aiproxy RealName "AI Proxy Service"
  dscl . -create /Users/_aiproxy IsHidden 1
  dscl . -create /Users/_aiproxy Password "*"
fi

install -d -m 0755 /etc/ai-proxy
install -d -m 0750 -o _aiproxy -g _aiproxy /var/lib/ai-proxy
install -d -m 0755 -o _aiproxy -g _aiproxy /var/log/ai-proxy

# Generate CA if absent (idempotent; preserves user CA on upgrade).
if [ ! -f /etc/ai-proxy/ca-cert.pem ] || [ ! -f /etc/ai-proxy/ca-key.pem ]; then
  /usr/local/bin/ai-proxy --generate-ca \
    --ca-cert /etc/ai-proxy/ca-cert.pem \
    --ca-key /etc/ai-proxy/ca-key.pem
  chmod 0644 /etc/ai-proxy/ca-cert.pem
  chmod 0640 /etc/ai-proxy/ca-key.pem
  chown _aiproxy:_aiproxy /etc/ai-proxy/ca-key.pem
fi

rollback_ca_trust() {
  # See linux postinstall.sh — same rationale: trusted CA without a running
  # interceptor is the failure mode that lets PII flow out unanonymized.
  # Exit code MUST reflect actual System keychain state, not intent: if the
  # delete fails or the cert is still trusted after, the caller (and the
  # operator) must see a distinct failure, not a successful-rollback signal.
  SHA=$(/usr/bin/openssl x509 -in /etc/ai-proxy/ca-cert.pem -noout -fingerprint -sha1 2>/dev/null \
        | sed 's/.*=//; s/://g' || true)
  if [ -z "$SHA" ]; then
    echo "ai-proxy: ROLLBACK FAILED — could not compute CA fingerprint to delete." >&2
    echo "ai-proxy:   manual remediation: 'sudo security delete-certificate -c \"ai-proxy CA\" /Library/Keychains/System.keychain'" >&2
    return 1
  fi
  /usr/bin/security delete-certificate -Z "$SHA" /Library/Keychains/System.keychain >/dev/null 2>&1
  # Verify the cert is actually gone — delete-certificate exit code alone is
  # insufficient (it returns 0 even when the cert never existed).
  if /usr/bin/security find-certificate -Z "$SHA" /Library/Keychains/System.keychain >/dev/null 2>&1; then
    echo "ai-proxy: ROLLBACK FAILED — CA $SHA still trusted in System keychain after delete attempt." >&2
    echo "ai-proxy:   PII MAY LEAK. Manual remediation required:" >&2
    echo "ai-proxy:     sudo security delete-certificate -Z $SHA /Library/Keychains/System.keychain" >&2
    return 1
  fi
  echo "ai-proxy: CA $SHA removed from System keychain (trust-without-interception averted)." >&2
  return 0
}

# Trust the CA in the System keychain.
if [ -f /etc/ai-proxy/ca-cert.pem ]; then
  /usr/bin/security add-trusted-cert -d -r trustRoot \
    -k /Library/Keychains/System.keychain \
    /etc/ai-proxy/ca-cert.pem >/dev/null 2>&1 || \
  /usr/bin/security add-trusted-cert -d -r trustRoot \
    -k /Library/Keychains/System.keychain \
    /etc/ai-proxy/ca-cert.pem
fi

# Bootstrap the LaunchDaemon. Capture stderr so the diagnostic survives the
# fallback path — silencing it would lose the actual error.
bootstrap_err=$(/bin/launchctl bootstrap system \
  /Library/LaunchDaemons/com.ai-anonymizing-proxy.plist 2>&1) || bootstrap_rc=$?
if [ "${bootstrap_rc:-0}" -ne 0 ]; then
  load_err=$(/bin/launchctl load -w \
    /Library/LaunchDaemons/com.ai-anonymizing-proxy.plist 2>&1) || load_rc=$?
  if [ "${load_rc:-0}" -ne 0 ]; then
    echo "ai-proxy postinstall: LaunchDaemon failed to load." >&2
    echo "ai-proxy:   launchctl bootstrap: $bootstrap_err" >&2
    echo "ai-proxy:   launchctl load -w:   $load_err" >&2
    echo "ai-proxy: rolling back CA trust to avoid trust-without-interception." >&2
    if ! rollback_ca_trust; then
      # Rollback itself failed — exit with a distinct code so the operator can
      # tell "service down, CA cleaned" (1) from "service down, CA STILL TRUSTED" (2).
      exit 2
    fi
    exit 1
  fi
fi

exit 0
